# laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php  /var/www/html/vachill.com/artisan queue:listen  --tries=3
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/vachill.com/storage/logs/supervisord.log
stopwaitsecs=30